# config/auditd/detectionlab.rules
-w /etc/passwd -p wa -k identity
-w /etc/group -p wa -k identity
-w /etc/shadow -p wa -k identity
-w /etc/sudoers -p wa -k scope
-w /etc/ssh/sshd_config -p wa -k sshd_config

# Cron & persistance
-w /etc/crontab -p wa -k cron
-w /etc/cron.d/ -p wa -k cron
-w /etc/cron.daily/ -p wa -k cron

# Web
-w /var/www/ -p wa -k webroot

# Commandes dangereuses
-a always,exit -F arch=b64 -S execve -F exe=/usr/bin/curl -k curl_exec
-a always,exit -F arch=b64 -S execve -F exe=/usr/bin/wget -k wget_exec
-a always,exit -F arch=b64 -S execve -F exe=/bin/nc -k netcat
-a always,exit -F arch=b64 -S execve -F exe=/usr/bin/python3 -k python_exec
